# deploy/rpi4/systemd/stream-watchdog.service
[Unit]
Description=Camera stream health check (one-shot)
After=camera-stream.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stream-watchdog.sh
StandardOutput=journal
StandardError=journal
